# Unix/installbuilder/service_scripts/omiserverd_sun10.xml
<?xml version="1.0" encoding="ISO-8859-1"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
        Copyright (c) Microsoft Corporation.  All rights reserved.
        Use is subject to license terms.

        ident   "@(#)omiserverd.xml   1.0     07/10/15"

        Service manifest for omiserver.
-->

<service_bundle type='manifest' name='OMI:omiserverd'>

  <service
     name='application/management/omiserverd'
     type='service'
     version='1'>
    
    <create_default_instance enabled='true' />
    
    <single_instance />
    
    <dependency name='network'
		grouping='require_all'
		restart_on='restart'
		type='service'>
      <service_fmri value='svc:/milestone/network' />
    </dependency>
    
    <dependency name='filesystem_local'
		grouping='require_all'
		restart_on='none'
		type='service'>
      <service_fmri value='svc:/system/filesystem/local' />
    </dependency>
    
    <exec_method
       type='method'
       name='start'
       exec='/opt/omi/bin/omiserverd %m'
       timeout_seconds='60' />
    
    <exec_method
       type='method'
       name='stop'
       exec='/opt/omi/bin/omiserverd %m'
       timeout_seconds='60' />
    
    <exec_method
       type='method'
       name='restart'
       exec='/opt/omi/bin/omiserverd %m'
       timeout_seconds='60' />
    
    <exec_method
       type='method'
       name='force-reload'
       exec='/opt/omi/bin/omiserverd %m'
       timeout_seconds='60' />
    
    <exec_method
       type='method'
       name='reload'
       exec='/opt/omi/bin/omiserverd %m'
       timeout_seconds='60' />
    
    <exec_method
       type='method'
       name='status'
       exec='/opt/omi/bin/omiserverd %m'
       timeout_seconds='60' />
    
    <property_group name='options' type='application'>
      <propval name='tcp_listen' type='boolean'
               value='false' />
    </property_group>
    
    <property_group name='startd' type='framework'>
      <!-- sub-process core dumps shouldn't restart session -->
      <propval name='ignore_error' type='astring' value='core,signal' />
    </property_group>
    
    <template>
      <common_name>
	<loctext xml:lang='C'>
          OMI Service
	</loctext>
      </common_name>
      
    </template>
  </service>
</service_bundle>
